# clock-browser.service
# systemd user service — launches Chromium in kiosk mode after the proxy starts.
#
# INSTALL (run as the clock user on the ArchLinux box):
#
#   cp clock-browser.service ~/.config/systemd/user/clock-browser.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now clock-browser.service
#
# REQUIRES: clock.service must already be installed and enabled.
#
# VIEW LOGS:
#   journalctl --user -u clock-browser.service -f

[Unit]
Description=Senior Clock kiosk browser
After=clock.service graphical-session.target
Requires=clock.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/chromium \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --no-first-run \
  --disable-session-crashed-bubble \
  --disable-restore-session-state \
  http://localhost:3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
